# etc/cron.daily/clean_docker_dangling_images
00 04 * * * root [ $(which docker) ] && docker rmi $(docker images -f "dangling=true" -q)
